# rootfs/etc/services.d/monitor/run
#!/command/with-contenv python3
"""s6 service run script: Monitor USB/IP devices and dependent apps.

Periodically checks that all expected USB/IP devices are still attached.
If a device is lost, attempts auto-reattach and notifies the user via
Home Assistant persistent notifications. Also monitors the health of
configured dependent apps and restarts them after successful reattach.
"""

import logging
import sys
import time

from usbip_lib.config import (
    get_app_config,
    send_ha_notification,
)
from usbip_lib.events import write_event
from usbip_lib.logging_setup import setup_logging
from usbip_lib.monitor import (
    attempt_reattach,
    check_dependent_app_health,
    find_missing_devices,
    is_on_cooldown,
    restart_dependent_apps,
    set_cooldown,
)
from usbip_lib.usbip import (
    parse_usbip_port,
    read_device_manifest,
)


def main() -> int:
    config = get_app_config()
    log_level = config.get("log_level", "info")
    logger = setup_logging(log_level, name="monitor")

    monitor_interval = config.get("monitor_interval", 30)
    reattach_retries = config.get("reattach_retries", 3)
    restart_retries = config.get("restart_retries", 3)

    logger.info(
        "Monitor service started (interval=%ds, reattach_retries=%d, restart_retries=%d)",
        monitor_interval, reattach_retries, restart_retries,
    )

    # Wait for the usbip service to do its initial attach
    logger.debug("Waiting 15s for initial device attachment to complete...")
    time.sleep(15)

    while True:
        try:
            # Re-read config each cycle to pick up changes from WebUI
            config = get_app_config()
            dependent_apps = config.get("dependent_apps", [])

            # Read manifest to know expected devices
            manifest = read_device_manifest()
            if not manifest:
                logger.debug("No device manifest found — nothing to monitor.")
                time.sleep(monitor_interval)
                continue

            # Get currently attached devices
            attached = parse_usbip_port()
            logger.debug(
                "Monitor check: %d expected, %d attached",
                len(manifest), len(attached),
            )

            # Detect missing devices
            missing = find_missing_devices(manifest, attached)

            if missing:
                logger.warning(
                    "%d device(s) missing: %s",
                    len(missing),
                    ", ".join(d.get("name", d["bus_id"]) for d in missing),
                )

                reattached_any = False

                for dev in missing:
                    bus_id = dev.get("bus_id", "")
                    name = dev.get("name", bus_id)
                    device_key = f"{dev.get('server', '')}:{bus_id}"

                    # Notify (with cooldown)
                    if not is_on_cooldown(device_key):
                        send_ha_notification(
                            "USB/IP: Device Lost",
                            f"{name} ({bus_id}) from {dev.get('server', '?')} is no longer attached. "
                            f"Attempting reattach...",
                            notification_type="device_lost",
                        )
                        write_event(
                            "device_lost",
                            f"{name} lost",
                            device=name,
                            server=dev.get("server", ""),
                        )
                        set_cooldown(device_key)

                    # Attempt reattach
                    ok = attempt_reattach(dev, reattach_retries, logger)
                    if ok:
                        reattached_any = True
                        send_ha_notification(
                            "USB/IP: Device Recovered",
                            f"{name} ({bus_id}) was successfully reattached.",
                            notification_type="device_recovered",
                        )
                        write_event(
                            "device_recovered",
                            f"{name} recovered",
                            device=name,
                            server=dev.get("server", ""),
                        )
                    else:
                        if not is_on_cooldown(f"fail:{device_key}"):
                            send_ha_notification(
                                "USB/IP: Reattach Failed",
                                f"Could not reattach {name} ({bus_id}) after "
                                f"{reattach_retries} attempts. Manual intervention required.",
                                notification_type="reattach_failed",
                            )
                            set_cooldown(f"fail:{device_key}")

                # Restart dependent apps after successful reattach
                if reattached_any and dependent_apps:
                    logger.info("Restarting dependent apps after device recovery...")
                    restart_dependent_apps(dependent_apps, restart_retries, logger)
            else:
                logger.debug("All %d expected device(s) are attached.", len(manifest))

            # Check dependent app health regardless of device status
            if dependent_apps:
                check_dependent_app_health(dependent_apps, restart_retries, logger)

        except (SystemExit, KeyboardInterrupt):
            logger.info("Monitor service shutting down.")
            raise
        except Exception:
            logger.exception("Monitor loop error")

        time.sleep(monitor_interval)


sys.exit(main())
