# rootfs/etc/services.d/webui/run
#!/command/with-contenv python3
"""s6 service run script: Launch the Flask WebUI."""

import logging
import os
import sys

logger = logging.getLogger("webui_run")
handler = logging.StreamHandler(sys.stdout)
handler.setFormatter(
	logging.Formatter(
		"%(asctime)s [%(name)s] %(levelname)s: %(message)s",
		datefmt="%H:%M:%S",
	)
)
logger.addHandler(handler)
logger.setLevel(logging.INFO)

logger.info("Starting WebUI service on port 8099...")

os.execvp("python3", ["python3", "/usr/local/bin/webui/app.py"])
