# rootfs/etc/services.d/usbip/finish
#!/command/with-contenv python3
"""s6 service finish script: Handle usbip service exit.

If the service crashed (non-zero exit code, excluding 256), halt the
entire container. Otherwise, s6 will automatically restart the service.
"""

import logging
import os
import sys

logger = logging.getLogger("usbip_finish")
handler = logging.StreamHandler(sys.stdout)
handler.setFormatter(
    logging.Formatter(
        "%(asctime)s [%(name)s] %(levelname)s: %(message)s",
        datefmt="%H:%M:%S",
    )
)
logger.addHandler(handler)
logger.setLevel(logging.INFO)

exit_code = int(sys.argv[1]) if len(sys.argv) > 1 else 0

if exit_code != 0 and exit_code != 256:
    logger.warning("usbip crashed (exit code %d), halting app", exit_code)
    os.execv("/run/s6/basedir/bin/halt", ["/run/s6/basedir/bin/halt"])
else:
    logger.info("usbip stopped, restarting...")
