# rootfs/etc/services.d/usbip/finish
#!/command/with-contenv python3
"""s6 service finish script: Handle usbip service exit.

The usbip service should be restarted by s6 when it exits. We avoid
halting the full container on transient usbip failures so the addon can
stay available and recover automatically.
"""

import logging
import sys

logger = logging.getLogger("usbip_finish")
handler = logging.StreamHandler(sys.stdout)
handler.setFormatter(
    logging.Formatter(
        "%(asctime)s [%(name)s] %(levelname)s: %(message)s",
        datefmt="%H:%M:%S",
    )
)
logger.addHandler(handler)
logger.setLevel(logging.INFO)

exit_code = int(sys.argv[1]) if len(sys.argv) > 1 else 0

if exit_code != 0 and exit_code != 256:
    logger.warning(
        "usbip exited unexpectedly (exit code %d), allowing s6 restart",
        exit_code,
    )
else:
    logger.info("usbip stopped, restarting...")
