# rootfs/etc/services.d/usbip/run
#!/command/with-contenv python3
"""s6 service run script: Attach USB/IP devices and keep the service alive.

Reads the device manifest created by init_devices.py, remounts sysfs,
attaches all devices, verifies attachment, and then sleeps forever
(s6 will call ``finish`` if this process exits).
"""

import signal
import sys
import time

from usbip_lib.config import get_app_config
from usbip_lib.events import write_event
from usbip_lib.logging_setup import setup_logging
from usbip_lib.usbip import (
    attach_all_from_manifest,
    parse_usbip_port,
    read_device_manifest,
    remount_sysfs,
    write_attached_devices_file,
)


def main() -> int:
    config = get_app_config()
    log_level = config.get("log_level", "info")
    logger = setup_logging(log_level, name="usbip_run")

    logger.info("Starting USB/IP devices attachment process.")

    # Read the device manifest
    manifest = read_device_manifest()
    if not manifest:
        logger.warning("No devices in manifest. Nothing to attach.")

    # Remount sysfs (required for USB/IP in containers)
    remount_sysfs()

    # Attach all devices
    succeeded, failed = attach_all_from_manifest(manifest)

    # Verify via usbip port
    time.sleep(1)
    attached = parse_usbip_port()
    attached_count = len(attached)

    # Write port numbers to file
    port_numbers = [d["port"] for d in attached]
    write_attached_devices_file(port_numbers)

    # Log results and write events
    if attached_count == 0:
        logger.error(
            "No USB/IP devices were attached. "
            "Is the USB/IP server online and device(s) USB attached?"
        )
        write_event("attach_fail", "No devices attached during startup")
    elif failed > 0:
        logger.warning(
            "%d device(s) attached, but some devices failed. "
            "Check configuration and server.",
            attached_count,
        )
        write_event("attach_ok", f"{attached_count} device(s) attached (partial)")
    else:
        logger.info("%d device(s) successfully attached.", attached_count)
        write_event("attach_ok", f"{attached_count} device(s) attached successfully")

    # Log each attached device
    for d in attached:
        logger.info(
            "Attached device on port %d: %s", d["port"], d.get("info", "")
        )

    logger.info("\U0001f7e2 USB/IP service has started")
    logger.info("\U0001f7e2 Keeping service alive.")

    # Keep the service running — s6 will restart if this exits
    try:
        while True:
            signal.pause()
    except KeyboardInterrupt:
        pass

    return 0


sys.exit(main())
