# rootfs/etc/services.d/usbip/run
#!/command/with-contenv bashio
# shellcheck disable=SC1008

bashio::log.info "Attaching USBIP devices"

# Run the attach device script
if /usr/local/bin/mount_devices; then
    bashio::log.info "USBIP devices attached successfully."
else
    bashio::log.error "Failed to attach USBIP devices. Is USBIP server online and device(s) usb attached?"
fi

# Keep the service running
exec sleep infinity
