# rootfs/etc/services.d/usbip/run
#!/command/with-contenv bashio
# shellcheck disable=SC1008
bashio::config.require 'log_level'
bashio::log.level "$(bashio::config 'log_level')"

bashio::log.info "Starting USB/IP devices attachment process."

# Run the attach device script
bashio::log.info "Running the mount_devices script."
mount_exit_code=0
/usr/local/bin/mount_devices || mount_exit_code=$?

# Check actual attachment status via usbip port (with brief delay for kernel to settle)
sleep 1
attached_devices=0
device_info_file="/tmp/attached_devices.txt"
rm -f "$device_info_file"

if ports_output=$(usbip port 2>/dev/null); then
    # Log each attached port
    while IFS= read -r port_line; do
        port_num=$(echo "$port_line" | sed -E 's/^Port ([0-9]+):.*/\1/')
        bashio::log.info "Attached device on port ${port_num}: $(echo "$port_line" | sed -E 's/^Port [0-9]+:\s*//')"
        echo "${port_num}" >> "$device_info_file"
        attached_devices=$((attached_devices + 1))
    done <<< "$(echo "$ports_output" | grep -E '^Port [0-9]+')"
fi

if [ "$attached_devices" -eq 0 ]; then
    bashio::log.error "No USB/IP devices were attached. Is the USB/IP server online and device(s) USB attached?"
    echo "{\"ts\":\"$(date -u +%Y-%m-%dT%H:%M:%SZ)\",\"type\":\"attach_fail\",\"device\":\"\",\"server\":\"\",\"detail\":\"No devices attached during startup\"}" >> /tmp/usbip_events.jsonl
elif [ "$mount_exit_code" -ne 0 ]; then
    bashio::log.warning "${attached_devices} device(s) attached, but some devices failed. Check configuration and server."
    echo "{\"ts\":\"$(date -u +%Y-%m-%dT%H:%M:%SZ)\",\"type\":\"attach_ok\",\"device\":\"\",\"server\":\"\",\"detail\":\"${attached_devices} device(s) attached (partial)\"}" >> /tmp/usbip_events.jsonl
else
    bashio::log.info "${attached_devices} device(s) successfully attached."
    echo "{\"ts\":\"$(date -u +%Y-%m-%dT%H:%M:%SZ)\",\"type\":\"attach_ok\",\"device\":\"\",\"server\":\"\",\"detail\":\"${attached_devices} device(s) attached successfully\"}" >> /tmp/usbip_events.jsonl
fi

bashio::log.info "🟢 USB/IP service has started"

# Keep the service running
bashio::log.info "🟢 Keeping service alive."
exec sleep infinity
